# src/wg-dashboard.service
[Unit]
After=netword.service

[Service]
WorkingDirectory=<your dashboard directory full path here>
ExecStart=/usr/bin/python3 <your dashboard directory full path here>dashboard.py
Restart=always


[Install]
WantedBy=default.target